# rootfs/etc/s6/olsrd/run
#!/bin/sh

olsrd -f /etc/olsrd/olsrd.conf -pidfile /tmp/olsrd.pid
sleep 5
pid=$(cat /tmp/olsrd.pid)
trap "kill $pid" TERM INT QUIT EXIT
while [ -e "/proc/$pid" ]; do
    sleep 1
done
exit 1
