# rootfs/etc/s6/olsrd/run
#!/bin/sh

olsrd

PID=$!

while kill -0 $PID > /dev/null 2>&1; do
  sleep 1
done
