# rootfs/etc/s6/syslogoutput/run
#!/bin/sh

cat /var/log/messages
exec tail -F /var/log/messages
