# rootfs/etc/s6/dnsmasq/run
#!/bin/sh

exec dnsmasq --keep-in-foreground
